# show/dmx-lights.service
# /etc/systemd/system/dmx-lights.service
#
#   sudo cp show/dmx-lights.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now dmx-lights
#   journalctl -u dmx-lights -f          # logs
#
[Unit]
Description=DMX stage lights
After=network.target dev-ttyUSB0.device
Wants=dev-ttyUSB0.device

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dmx-ai/show
ExecStart=/usr/bin/python3 /home/pi/dmx-ai/show/app.py --http-port 80
Restart=always
RestartSec=3
# Let it bind :80 without root
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
